# systemd/mdadm-last-resort@.service
[Unit]
Description=Activate md array %I even though degraded
DefaultDependencies=no
ConditionPathExists=!/sys/devices/virtual/block/%I/md/sync_action

[Service]
Type=oneshot
ExecStart=BINDIR/mdadm --run /dev/%I
